# backup/rdiffbackup/rdiffbackup-init.conf
[program:init]
autorestart = false
startsecs = 0
user = root
command = /root/init.sh